# void/sv/lte/finish
#!/bin/sh
# Stop LTE
/usr/comma/lte/lte.sh stop 2>/dev/null || true
