# void/sv/sound/run
#!/bin/sh
exec /usr/comma/sound/sound_init.sh
