# void/sv/sound/run
#!/bin/sh
# One-shot service - run once then sleep forever
/usr/comma/sound/sound_init.sh
exec sleep infinity
