# void/sv/irsc_util/run
#!/bin/sh
# irsc_util - Security config (runs once at boot)
exec 2>&1

sv check dbus > /dev/null || exit 1

exec /usr/bin/irsc_util /etc/sec_config
