# void/sv/avahi-ssh-publish/run
#!/bin/sh
exec /usr/comma/avahi-publish-service.sh
