# userspace/rootfs/etc/systemd/system/qseecomd.path
# Copyright (c) 2018 Qualcomm Technologies, Inc.
# All Rights Reserved.
# Confidential and Proprietary - Qualcomm Technologies, Inc.
#

[Unit]
Description=Path watching for qseecomd.service
DefaultDependencies=no
Conflicts=shutdown.target
Before=paths.target shutdown.target systemd-udev-trigger.service

[Path]
PathModified=/dev/qseecom

[Install]
WantedBy=multi-user.target
