# void/sv/pdmapper/run
#!/bin/sh
# pdmapper - Protection Domain mapper service
exec 2>&1

sv check init_mss > /dev/null || exit 1

# One-shot service - run once then sleep forever
/etc/initscripts/start_pdmappersvc start
exec sleep infinity
