# void/sv/tftp_server/run
#!/bin/sh
# tftp_server - Qualcomm TFTP server for firmware
exec 2>&1

sv check rmt_storage > /dev/null || exit 1

exec /usr/bin/tftp_server
