# void/sv/leprop/run
#!/bin/sh
# leprop - Android-style property service
exec 2>&1

exec /sbin/leprop-service
