# void/sv/leprop/run
#!/bin/sh
# leprop - Android-style property service
exec 2>&1

# Create /dev/socket directory for leprop and logd sockets
mkdir -p /dev/socket
chmod 777 /dev/socket

exec /usr/bin/leprop-service
