# void/sv/gpio/run
#!/bin/sh
# One-shot service - run once then sleep forever
/usr/comma/gpio.sh
exec sleep infinity
