# void/sv/magic/run
#!/bin/sh
# magic - GPU display service (raylib on Adreno GPU)
exec 2>&1

if [ ! -e /dev/dri/card0 ]; then
    echo "Waiting for /dev/dri/card0..."
    exit 1
fi

# Suppress Qualcomm display driver debug spam
mkdir -p /data/misc/display
echo 0 > /data/misc/display/sdm_dbg_cfg.txt 2>/dev/null || true
echo 0 > /data/misc/display/gbm_dbg_cfg.txt 2>/dev/null || true

exec chpst -u comma:comma:gpu:video /usr/local/venv/bin/python -u /usr/comma/magic.py
