# void/sv/magic/run
#!/bin/sh
# magic - GPU display magic (setup permissions and run python)
exec 2>&1

# Wait for GPU device to exist
if [ ! -e /dev/dri/card0 ]; then
    echo "Waiting for /dev/dri/card0..."
    exit 1
fi

# Setup GPU device permissions
chgrp gpu /dev/ion /dev/kgsl-3d0 /dev/dri/card0 /dev/dri/renderD128 2>/dev/null || true
chmod 660 /dev/ion /dev/kgsl-3d0 /dev/dri/card0 /dev/dri/renderD128 2>/dev/null || true

# Verify comma user can access GPU before starting python
if ! su -s /bin/sh comma -c "test -r /dev/dri/card0 && test -w /dev/dri/card0"; then
    echo "GPU not accessible by comma user"
    exit 1
fi

# Setup display debug config
mkdir -p /data/misc/display
echo 0 > /data/misc/display/sdm_dbg_cfg.txt 2>/dev/null || true
echo 0 > /data/misc/display/gbm_dbg_cfg.txt 2>/dev/null || true

# Run as comma user
exec chpst -u comma /usr/local/venv/bin/python -u /usr/comma/magic.py
