# void/sv/init_mss/run
#!/bin/sh
# init_mss - Modem subsystem init (runs once)
exec 2>&1

sv check rmt_storage > /dev/null || exit 1

exec /usr/sbin/init_mss
