# void/sv/screen_calibration/run
#!/bin/sh
exec /usr/comma/screen_calibration.py
