# void/sv/varwatch/run
#!/bin/sh
. /etc/profile
exec /usr/local/venv/bin/python -u /usr/comma/varwatch.py
