# void/files/50-firmware.rules
# Firmware loading helper for eudev - responds immediately to prevent 60s timeouts
SUBSYSTEM=="firmware", ACTION=="add", RUN+="/lib/udev/firmware-helper.sh"
