# void/sv/cdsp/run
#!/bin/sh
# cdsp - Compute DSP init (runs once)
exec 2>&1

# Run the init script
/usr/bin/cdsp.sh

# Stay up so runit doesn't restart
exec sleep infinity
