# userspace/files/90-eth-mac.rules
ACTION=="add", SUBSYSTEM=="net", DRIVERS=="r8152", RUN+="/usr/comma/set-eth-mac.sh %k"
