# userspace/files/adb-param-watcher.path
[Path]
PathModified=/data/params/d/AdbEnabled

[Install]
WantedBy=multi-user.target
