# userspace/files/agnos-debug.service
[Unit]
After=local-fs.target

[Service]
Type=oneshot
User=comma
ExecStart=/usr/comma/debug.sh

[Install]
WantedBy=multi-user.target
