# userspace/files/init-qcom.service
[Unit]
Description=QCOM init
After=sysinit.target

[Service]
Type=oneshot
ExecStart=/usr/comma/init.qcom.sh

[Install]
WantedBy=multi-user.target
