# userspace/files/set_time.service
[Unit]

[Service]
Type=oneshot
ExecStart=/usr/local/pyenv/shims/python -u /usr/comma/set_time.py

[Install]
WantedBy=multi-user.target
