# userspace/files/sound.service
[Unit]
Description=Sound

[Service]
Type=oneshot
ExecStart=/usr/comma/sound/sound_init.sh

[Install]
WantedBy=multi-user.target
