# userspace/files/ssh-param-watcher.service
[Unit]
Description=SSH param watcher

[Service]
Type=oneshot
ExecStart=/usr/comma/set_ssh.sh

[Install]
WantedBy=rescue.target
