# userspace/files/ssh-param-watcher.path
[Path]
PathModified=/data/params/d/SshEnabled

[Install]
WantedBy=multi-user.target
