# userspace/files/weston.service
[Unit]
Description=Weston

[Service]
Type=simple
Restart=always
Environment="XDG_RUNTIME_DIR=/var/tmp/weston"

# If we start weston immediately it just exits without an error
# TODO: Figure out what we need to wait for (drm?)
# run modetest for each panel. unclear why weston doesn't take care of it though
ExecStart=/bin/bash -c "sleep 4 && \
                        /usr/comma/modetest -M msm_drm -s 26@109:1080x2160-60 && \
                        /usr/comma/modetest -M msm_drm -s 26@111:1080x2160-60 && \
                        sleep 1 && \
                        mkdir -p $XDG_RUNTIME_DIR && \
                        chown -R comma: $XDG_RUNTIME_DIR && \
                        mkdir -p /data/misc/display || true && \
                        echo 0 > /data/misc/display/sdm_dbg_cfg.txt || true && \
                        echo 0 > /data/misc/display/gbm_dbg_cfg.txt || true && \
                        /usr/bin/weston --idle-time=0 --tty=1 --config=/usr/comma/weston.ini"
KillSignal=SIGKILL

[Install]
WantedBy=multi-user.target
