# void/sv/serial-hostname/run
#!/bin/sh
exec /usr/comma/set-hostname.sh
