# void/sv/serial-hostname/run
#!/bin/sh
# One-shot service - run once then sleep forever
/usr/comma/set-hostname.sh
exec sleep infinity
