# void/sv/power_drop_monitor/finish
#!/bin/sh
# Delay restart to avoid spamming console on failure
sleep 5
